# docs/legacy/benchmarks/jff-arch-config/modules/services/systemd-user/units/jff-secrets-check.timer
[Unit]
Description=JFF Secret Service Health Check Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Unit=jff-secrets-check.service
Persistent=true

[Install]
WantedBy=timers.target
