# docs/legacy/benchmarks/jff-arch-config/modules/services/systemd-user/units/jff-secrets-check.service
[Unit]
Description=JFF Secret Service Health Check
After=graphical-session.target

[Service]
Type=oneshot
WorkingDirectory=%h/Documents/jff-arch-config
ExecStart=%h/Documents/jff-arch-config/scripts/check-secrets-service.sh
