# docs/legacy/benchmarks/jff-arch-config/modules/services/systemd-user/units/jff-update-dry-run.timer
[Unit]
Description=JFF Arch Update Dry-Run Timer

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=true
Unit=jff-update-dry-run.service

[Install]
WantedBy=timers.target
