# modules/audit/configs/audit-rules.conf
# Iron Security Audit Rules
# Deployed by iron-arch audit module

# Self-auditing
-w /var/log/audit/ -k auditlog
-w /etc/audit/ -p wa -k auditconfig
-w /etc/libaudit.conf -p wa -k auditconfig

# Authentication events
-w /etc/pam.d/ -p wa -k auth_pam
-w /etc/login.defs -p wa -k auth_logins
-w /etc/securetty -p wa -k auth_securetty
-w /var/log/faillog -p wa -k auth_faillog
-w /var/log/lastlog -p wa -k auth_lastlog

# System administration
-w /etc/localtime -p wa -k time_change
-a always,exit -F arch=b64 -S adjtimex -S settimeofday -k time_change
-w /etc/hostname -p wa -k hostname_change
-w /etc/cron.d/ -p wa -k cron
-w /etc/crontab -p wa -k cron

# Network configuration
-w /etc/hosts -p wa -k network_hosts
-w /etc/resolv.conf -p wa -k network_dns
-a always,exit -F arch=b64 -S sethostname -S setdomainname -k network_hostname

# Kernel modules
-w /sbin/insmod -p x -k kernel_modules
-w /sbin/rmmod -p x -k kernel_modules
-w /sbin/modprobe -p x -k kernel_modules

# File system
-a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -F auid>=1000 -F auid!=4294967295 -k file_perm
-a always,exit -F arch=b64 -S chown -S fchown -S fchownat -F auid>=1000 -F auid!=4294967295 -k file_owner

# Privilege escalation
-w /bin/su -p x -k privilege_escalation
-w /usr/bin/sudo -p x -k privilege_escalation
-w /etc/sudoers -p wa -k privilege_escalation
-w /etc/sudoers.d/ -p wa -k privilege_escalation

# Package management
-w /usr/bin/pacman -p x -k package_management
-w /etc/pacman.conf -p wa -k package_management

# SSH
-w /etc/ssh/sshd_config -p wa -k sshd_config
-w /etc/ssh/sshd_config.d/ -p wa -k sshd_config
